# bot-rw.service.sample
[Unit]
Wants = network.target
After = multi-user.target
Description = Discord bot service

[Service]
Type = simple
Restart = always
Restart = on-failure

Group = user
User = user
WorkingDirectory = /home/user/bot-rw
ExecStart = /usr/bin/python3 /home/user/bot-rw/run.py

[Install]
WantedBy = multi-user.target
